# 99-usb-photo-upload.rules
# udev rule to trigger USB photo upload when a USB block device is added
# This rule is installed to /etc/udev/rules.d/ by the install script

# Trigger on any block device (partition) being added
# The service itself filters to only process USB devices
SUBSYSTEM=="block", ACTION=="add", RUN+="/usr/bin/systemctl try-restart usb-photo-upload.service"
