# usb-photo-upload.service
[Unit]
Description=USB Photo Auto Upload to Synology NAS
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=@INSTALL_DIR@
ExecStart=@INSTALL_DIR@/venv/bin/python3 @INSTALL_DIR@/usb_photo_upload.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
